# installer/mac/scout-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>ai.scout.bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Scout Bridge.app/Contents/Resources/node/bin/node</string>
        <string>/Applications/Scout Bridge.app/Contents/Resources/bridge/dist/index.js</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

<<<<<<< HEAD
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>

    <key>StandardErrorPath</key>
    <string>/tmp/scout-bridge-error.log</string>

</dict>
</plist>
=======
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/bridge.log</string>

    <key>StandardErrorPath</key>
